# run/supervisord.conf
[supervisord]
nodaemon=true

[program:main]
command=python /app/main.py
stdout_logfile=/var/log/main_cvr.log
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
startsecs=0
autorestart=true
